# callstats-grafana/ats/scripts/callstats_update.cron
# callstats-grafana ETL — обновление статистики каждые 5 минут
# Устанавливается в /etc/cron.d/callstats_update на АТС
*/5 * * * * root /usr/local/bin/callstats_update.sh >>/var/log/callstats_update.log 2>&1
